# users/bruh/src/config/gnome-autostart/RemoveHyprlandLast.desktop
[Desktop Entry]
Type=Application
Name=RemoveHyprlandLast
Exec=rm /home/bruh/.cache/.hyprland_last
