# users/bruh/src/config/gnome-autostart/StartXWayland.desktop
[Desktop Entry]
Type=Application
Name=StartXWayland
Exec=xset -q > /dev/null
