# users/bruh/src/config/gnome-autostart/UnhideTitlebars.desktop
[Desktop Entry]
Type=Application
Name=UnhideTitlebars
Exec=dconf write /org/gnome/shell/extensions/unite/hide-window-titlebars "'never'"
