# users/bruh/src/config/gnome-autostart/UnmaskGnomePortal.desktop
[Desktop Entry]
Type=Application
Name=UnmaskGnomePortal
Exec=systemctl --user unmask xdg-desktop-portal-gnome.service
